# vagrant/modules/collectd/files/collectd-upstart.conf
description "collectd"
start on runlevel [235]
stop on runlevel [016]
respawn
exec /opt/collectd/sbin/collectd -f
