# vagrant/modules/shorten/files/shorten-upstart.conf
description "python url shortener"
start on runlevel [235]
stop on runlevel [016]
respawn
exec gunicorn url_shortener:app