# vagrant/modules/riemann/files/riemann-dash-upstart.conf
description "riemann dashboard"
start on runlevel [235]
stop on runlevel [016]
respawn
exec riemann-dash /etc/riemann-dash.conf